# scripts/com.sergey.investments-monitor.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sergey.investments-monitor.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/sergeypochikovskiy/AI_workspace/Projects/4. Investments Monitor/scripts/refresh.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sergeypochikovskiy/.logs/investments-monitor-refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sergeypochikovskiy/.logs/investments-monitor-refresh-error.log</string>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
